# config/zpa-report.service
[Unit]
Description=ZPA Status Mini-SIEM — daily report generator
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/zpa-siem/venv/bin/python3 /opt/zpa-siem/report_generator.py
WorkingDirectory=/opt/zpa-siem
